# i3/X/turnon-laptop-screen-if-last/99-turnon-laptop-screen-if-last.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/fanz/.local/bin/turnon-laptop-screen-if-last.sh"
